# lock-activated.target
[Unit]
Conflicts=unlock-activated.target
Description=This is activated when the session is locked
